# packaging_ubuntu/postrm
#!/bin/bash
# This `DEBIAN/postrm` script is run post-uninstallation

set -e

systemctl stop kamut-node    || true
systemctl disable kamut-node || true
rm -f /etc/systemd/system/kamut-node.service
